# app/conf/supervisor_flask.ini
[supervisord]
nodaemon=true

[program:flaskweb]
command=gunicorn --worker-class quart.worker.GunicornWorker -w 4 -b 0.0.0.0:5000 app:app
directory=/app/web
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/app/logs/quart_out.log
stderr_logfile=/app/logs/quart_err.log
